# lib/systemd/system/vcd-sub.service
[Unit]
Description=Vitalinux Consumer Daemon — Sub-instance Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/vcd-sub-manager
Restart=always
RestartSec=30
RuntimeDirectory=vcd

[Install]
WantedBy=multi-user.target
