# debian/postinst
#!/bin/sh
set -e

echo "[val] Ejecutando postinst..."

mkdir -p /etc/val/val.conf.d
mkdir -p /etc/val/hooks.d
mkdir -p /etc/val/val.sub
mkdir -p /var/lib/val
mkdir -p /var/lib/val/sub

if [ ! -f /var/lib/val/version ]; then
    echo "0" > /var/lib/val/version
fi

systemctl daemon-reload || true
systemctl enable val.service || true

if [ "$1" = "configure" ]; then
    if [ -n "${2:-}" ]; then
        systemctl restart val.service || true
    else
        systemctl start val.service || true
    fi
fi

echo "[val] Instalacion completada."
exit 0
