# debian/postinst
#!/bin/sh
set -e

echo "[vcd] Ejecutando postinst..."

mkdir -p /etc/vcd/vcd.conf.d
mkdir -p /etc/vcd/hooks.d
mkdir -p /var/lib/vcd

if [ ! -f /var/lib/vcd/version ]; then
    echo "0" > /var/lib/vcd/version
fi

systemctl daemon-reload || true
systemctl enable vcd.service || true

if [ "$1" = "configure" ]; then
    if [ -n "${2:-}" ]; then
        systemctl restart vcd.service || true
    else
        systemctl start vcd.service || true
    fi
fi

echo "[vcd] Instalacion completada."
